# lib/flipstone-deployment/capistrano/rails/upstartjob.conf.erb
# <%= application %> upstart job file
# /etc/init/<%= application %>.conf
# This configuration manages <%= application %>

description "start <%= application %>"

start on runlevel [23]
stop on runlevel [01456]

chdir <%=current_path%>

exec start-stop-daemon -c <%=user%> --pidfile <%=shared_path%>/pids/unicorn.pid --start --startas /usr/local/rvm/bin/<%= application %>_unicorn -- -E <%= rails_env %> -c <%= shared_path %>/system/unicorn.conf

respawn

