# lib/flipstone-deployment/capistrano/rails/upstartjob.conf.erb
# <%= application %> upstart job file
# /etc/init/<%= application %>.conf
# This configuration manages <%= application %>

description "start <%= application %>"

start on runlevel [23]
stop on runlevel [01456]

chdir <%=current_path%>

exec sudo su -c 'rvm 1.9.2 exec bundle exec <%= unicorn[:process] || 'unicorn' %> -E <%=rails_env%> -c <%=shared_path%>/system/unicorn.conf' <%=user%>

respawn

